# etc/xinetd.d/ves-smtp
service ves-smtp
{
	disable	= no
	type		= UNLISTED
	port		= 7125
	socket_type     = stream
	wait            = no
	user            = mail
	server          = /usr/local/bin/vesmail
	server_args     = smtp --cert=/etc/exim.crt --pkey=/etc/exim.key --ca=/etc/exim.crt
	log_on_failure  += USERID
	instances	= 100
	per_source	= 25
}
