# etc/xinetd.d/ves-imap
service ves-imap
{
	disable	= no
	type		= UNLISTED
	port		= 7143
	socket_type     = stream
	wait            = no
	user            = mail
	server          = /usr/local/bin/vesmail
	server_args     = imap
	log_on_failure  += USERID
	instances	= 100
	per_source	= 25
}
