# etc/systemd/system/vesmail.service
[Unit]
Description=VESmail Proxy

[Service]
User=vesmail
Group=vesmail
ExecStart=/usr/bin/vesmail daemon -GG
